# services/drop.service
[Unit]
Description=Protected URL download service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/usr/bin/python3 /opt/url-drop/drop.py
Environment=PYTHONUNBUFFERED=1
Environment=DROP_SOCKET=/run/url-drop/url-drop.sock
Environment=DROP_DIR=/var/lib/url-drop/files
Environment=DROP_MAX_BYTES=21474836480
Environment=DROP_RESERVE_BYTES=5368709120
RuntimeDirectory=url-drop
RuntimeDirectoryMode=0750
StateDirectory=url-drop
StateDirectoryMode=0750
UMask=0077
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
RemoveIPC=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressAllow=127.0.0.53
IPAddressDeny=0.0.0.0/8
IPAddressDeny=127.0.0.0/8
IPAddressDeny=10.0.0.0/8
IPAddressDeny=100.64.0.0/10
IPAddressDeny=169.254.0.0/16
IPAddressDeny=172.16.0.0/12
IPAddressDeny=192.0.0.0/24
IPAddressDeny=192.0.2.0/24
IPAddressDeny=192.168.0.0/16
IPAddressDeny=198.18.0.0/15
IPAddressDeny=198.51.100.0/24
IPAddressDeny=203.0.113.0/24
IPAddressDeny=224.0.0.0/4
IPAddressDeny=240.0.0.0/4
IPAddressDeny=::1
IPAddressDeny=100::/64
IPAddressDeny=2001:db8::/32
IPAddressDeny=fc00::/7
IPAddressDeny=fe80::/10
IPAddressDeny=ff00::/8

[Install]
WantedBy=multi-user.target
